# root/etc/services.d/sshfs/run
#!/usr/bin/with-contenv bash

_term() {
	echo "Caught SIGTERM for SSHFS!"
	fusermount3 -uz ${SSHFS_MNT}
}

trap _term SIGTERM

exec s6-setuidgid abc sshfs -f ${SSHFS_HOST} ${SSHFS_MNT} ${SSHFS_OPTS}
